# scripts/pihole-do/daily-reboot.timer
[Unit]
Description=Daily reboot

[Timer]
OnCalendar=*-*-* 3:30:00
Persistent=true

[Install]
WantedBy=timers.target
